# deploy/investmentology-daily-analyze.service
[Unit]
Description=Investmentology Daily Watchlist Re-Analysis
After=network.target postgresql.service

[Service]
Type=oneshot
WorkingDirectory=/home/investmentology
ExecStart=/usr/bin/python3 -m investmentology cron daily-watchlist-analyze --limit 30
Environment=PYTHONPATH=/home/investmentology/src
StandardOutput=journal
StandardError=journal
TimeoutStartSec=7200
