# deploy/investmentology-monitor.service
[Unit]
Description=Investmentology Daily Post-Market Monitor
After=network.target postgresql.service

[Service]
Type=oneshot
WorkingDirectory=/home/investmentology
ExecStart=/usr/bin/python3 -m investmentology cron daily-monitor
Environment=PYTHONPATH=/home/investmentology/src
StandardOutput=journal
StandardError=journal
TimeoutStartSec=1800
