# deploy/investmentology-monitor.timer
[Unit]
Description=Investmentology Daily Post-Market Monitor Timer

[Timer]
OnCalendar=Mon..Fri 16:15 America/New_York
Persistent=true

[Install]
WantedBy=timers.target
